# src/server/timed-qt5.service
[Unit]
Description=Time Daemon
Requires=dbus.socket

[Service]
Type=notify
ExecStart=/usr/bin/timed-qt5 --systemd
Restart=always

[Install]
WantedBy=pre-user-session.target
