# src/server/timed.service
[Unit]
Description=Time Daemon
After=syslog.target

[Service]
ExecStart=/usr/bin/timed

[Install]
WantedBy=multi-user.target
