# src/server/timed.service
[Unit]
Description=Time Daemon
Requires=dbus.socket

[Service]
Type=notify
ExecStart=/usr/bin/timed --systemd
Restart=always

[Install]
WantedBy=pre-user-session.target
